# mongo-data-rs0/mongo-data-rs0.conf
[program:mongo-data-rs0]
command         = mongod --dbpath /data/rs0 --rest --bind_ip 0.0.0.0 --port 27017 --replSet rs0
autostart       = true
autorestart     = true
exitcodes       = 2
numprocs        = 1
startsecs       = 5
user            = root
