# rpi-spider-data/mysql_setup.conf
[program:mysql-setup]
command     = /root/mysql_setup.sh
autostart   = true
autorestart = false
exitcodes   = 0
numprocs    = 1
startsecs   = 0
user        = root
