# spider-data/beanstalkd.conf
[program:beanstalkd]
command     = beanstalkd -l 0.0.0.0 -p 11300
autostart   = true
autorestart = true
exitcodes   = 2
user        = nobody
