# mongo-monitoring/mongomms.conf
[program:MongoMMS]
command     = /usr/bin/mongodb-mms-monitoring-agent -conf /etc/mongodb-mms/monitoring-agent.config
autostart   = true
autorestart = true
user        = mongodb-mms-agent
